# crontab
*/5 * * * * /home/docker/cronjobs/date.sh